# blackbox-server/blackbox-server_1.0.0/DEBIAN/postinst
#!/bin/bash
set -e

# Enable and start systemd service
if [ -d /run/systemd/system ]; then
    systemctl daemon-reload
    systemctl enable blackbox-server.service || true
    systemctl start blackbox-server.service || true
fi

exit 0
