# blackbox-server/debian/blackbox-server.service
[Unit]
Description=Blackbox GPU VRAM Monitoring Server
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/blackbox-server 6767
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target


